# deploy/systemd/netdata.service.d/limits.conf
# RU — cap netdata so monitoring (Phase 1 addition) can't load the 2-vCPU box
# that all VPN country relays cascade through. Netdata is sacrificeable → +OOM.
# Install: /etc/systemd/system/netdata.service.d/limits.conf ; daemon-reload + restart.
# Pair with /etc/netdata/netdata.conf [global] update every = 2 (halves collection).
# Rollback: rm this file && systemctl daemon-reload && systemctl restart netdata.
[Service]
CPUQuota=20%
MemoryHigh=180M
MemoryMax=240M
OOMScoreAdjust=300
