# deploy/systemd/ip-watchdog.service
[Unit]
Description=ip-watchdog DNS failover
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=bozhenkas
EnvironmentFile=/etc/ip-watchdog/watchdog.env
ExecStart=/usr/bin/python3 /opt/ip-watchdog/watchdog.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ip-watchdog
