# deploy/systemd/sub-updater.timer
[Unit]
Description=Remnawave Hydra subscription updater timer (every 30 min)
Requires=sub-updater.service

[Timer]
# Run shortly after boot, then every 30 minutes (matches old sleep 1800).
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=60s
Persistent=true

[Install]
WantedBy=timers.target
