# deploy/systemd/sub-updater.service
[Unit]
Description=Remnawave Hydra subscription updater
After=network.target docker.service remnawave.service

[Service]
Type=simple
Environment=REMNA_PUBLIC_PORT=443
Environment=REMNA_NGINX_CONF=/etc/nginx/sites-enabled/ru.goida.fun
ExecStart=/bin/bash -lc 'while true; do /usr/bin/python3 /root/goida-scripts/sync_remna_hydra.py || true; sleep 1800; done'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
