# deploy/systemd/telegram-proxy-tunnel.service
[Unit]
Description=Telegram proxy tunnel to FIN
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ssh -N -T \
    -i /root/.ssh/telegram_proxy_ed25519 \
    -o BatchMode=yes \
    -o ExitOnForwardFailure=yes \
    -o ServerAliveInterval=30 \
    -o ServerAliveCountMax=3 \
    -o StrictHostKeyChecking=yes \
    -p 17904 \
    -L 127.0.0.1:8888:127.0.0.1:8888 \
    root@77.110.108.57
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
