# deploy/systemd/vpn-bot.service.d/cluster.conf
[Service]
Environment=CLUSTER_PLAIN_PATH=/etc/goida/cluster.json
Environment=CLUSTER_CACHE_PATH=/var/lib/goida/cluster-cache.json
# CLUSTER_HMAC_SECRET задаётся через secrets.env (gitignored) или Ansible vault
# EnvironmentFile=/etc/goida/secrets.env
